# output/sync.ini
[supervisord]
nodaemon=true

[program:mongosync]  
command=/opt/mongosync/output/mongosync -c /opt/mongosync/output/mongosync.conf